# tensorflow/22.pyspark/build/supervisord.conf
[supervisord]

[program:pyspark-master]
command = bash -c "/opt/spark/bin/spark-class org.apache.spark.deploy.master.Master"
stopasgroup = true
autorestart = true

[inet_http_server]
port = 0.0.0.0:9000