# piping/2.luigi-crawler-sentiment-elasticsearch/supervisord.conf
[supervisord]

[program:luigi]
command = luigid
stopasgroup = true
autorestart = true

[inet_http_server]
port = 0.0.0.0:9000
